# packaging/droplet/systemd/celiums-proxy.service
[Unit]
Description=Celiums reverse proxy (Caddy: TLS + path routing → engine + dashboard)
Documentation=https://github.com/terrizoaguimor/celiums-memory
After=network-online.target celiums-dashboard.service celiums-engine.service
Wants=network-online.target
Requires=celiums-dashboard.service celiums-engine.service

[Service]
Type=notify
User=celiums
Group=celiums
ExecStart=/usr/bin/caddy run --config /etc/celiums/Caddyfile --adapter caddyfile
ExecReload=/usr/bin/caddy reload --config /etc/celiums/Caddyfile --adapter caddyfile --force
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s

# Caddy now terminates TLS for the public listeners (80 + 443). Both ports
# are privileged, so we grant the binary the bind-service capability — no
# need to run the whole process as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Caddy needs a working dir for its certificate cache (Let's Encrypt mode).
ReadWritePaths=/var/lib/celiums /var/lib/caddy
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=celiums-proxy

[Install]
WantedBy=multi-user.target
